# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stopasgroup=true

[program:crond]
command=/usr/sbin/crond -f
stopasgroup=true
